# home/Library/LaunchAgents/com.jeremybanka.helix-theme.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.jeremybanka.helix-theme</string>
	<key>ProgramArguments</key>
	<array>
		<string>/opt/homebrew/bin/bun</string>
		<string>/Users/jem/dotfiles/scripts/sync-helix-theme.ts</string>
		<string>--watch</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/Users/jem/Library/Logs/helix-theme-sync.err.log</string>
	<key>StandardOutPath</key>
	<string>/Users/jem/Library/Logs/helix-theme-sync.out.log</string>
</dict>
</plist>
