# home/Library/LaunchAgents/com.jeremybanka.helix-theme.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.jeremybanka.helix-theme</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/zsh</string>
		<string>-lc</string>
		<string>exec /opt/homebrew/bin/bun "$HOME/dotfiles/home/Library/LaunchAgents/com.jeremybanka.helix-theme/theme-switcher.ts" --watch</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/tmp/com.jeremybanka.helix-theme.err.log</string>
	<key>StandardOutPath</key>
	<string>/tmp/com.jeremybanka.helix-theme.out.log</string>
</dict>
</plist>
